# jarvis_backup.timer
[Unit]
Description=Run Jarvis local memory backup periodically

[Timer]
OnCalendar=hourly
Persistent=true
Unit=jarvis_backup.service

[Install]
WantedBy=timers.target
