# molt-evo-sync.service
[Unit]
Description=molt_evo knowledge sync
After=network.target

[Service]
Type=oneshot
Environment=MOLT_EVO_WORKSPACE=/root/.openclaw/workspace
Environment=MOLT_EVO_REPO_DIR=/github/molt_evo
WorkingDirectory=/root/.openclaw/workspace
ExecStart=/bin/bash -lc 'python3 $MOLT_EVO_REPO_DIR/molt_evo_daily_summary.py'
