# molt-evo-sync.service
[Unit]
Description=molt_evo knowledge sync
After=network.target

[Service]
Type=oneshot
WorkingDirectory=${MOLT_EVO_WORKSPACE:-/root/.openclaw/workspace}
ExecStart=/bin/bash -lc 'python3 ${MOLT_EVO_WORKSPACE:-/root/.openclaw/workspace}/molt_evo_daily_summary.py'
