# molt-evo-backup.service
[Unit]
Description=molt_evo local memory backup
After=network.target

[Service]
Type=oneshot
WorkingDirectory=${MOLT_EVO_WORKSPACE:-/root/.openclaw/workspace}
ExecStart=/bin/bash ${MOLT_EVO_WORKSPACE:-/root/.openclaw/workspace}/molt_evo_backup.sh
