# molt-evo-backup.timer
[Unit]
Description=Run molt_evo local memory backup periodically

[Timer]
OnCalendar=hourly
Persistent=true
Unit=molt-evo-backup.service

[Install]
WantedBy=timers.target
