# release-excluded/compat/legacy-jarvis/jarvis_backup.service
[Unit]
Description=Jarvis local memory backup
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/root/.openclaw/workspace
ExecStart=/bin/bash /root/.openclaw/workspace/jarvis_backup.sh
